# systemd/trade-cache-refresh.timer
[Unit]
Description=Refresh trade dashboard cache every 5 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
